# macos-ui/packaging/launchd/com.github.ms-sso-openconnect-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this LaunchAgent -->
    <key>Label</key>
    <string>com.github.ms-sso-openconnect-ui</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/MS SSO OpenConnect.app/Contents/MacOS/ms-sso-openconnect-ui</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep it running (restart if it crashes) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- This is a GUI app, allow it to interact with desktop -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Set environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PLAYWRIGHT_BROWSERS_PATH</key>
        <string>~/Library/Caches/ms-playwright</string>
    </dict>

    <!-- Log files -->
    <key>StandardOutPath</key>
    <string>~/Library/Logs/ms-sso-openconnect-ui.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/ms-sso-openconnect-ui.error.log</string>

    <!-- Delay startup slightly to let desktop fully initialize -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
